# supervisord/nginx.conf
[program:nginx]
command=nginx
numprocs=1
autostart=true
autorestart=true