# supervisord/php5-fpm.conf
[program:php5-fpm]
command=php5-fpm
numprocs=1
autostart=true
autorestart=true